# deploy/supervisor/web-smk1kolaka-queue.conf
; Supervisor config for Laravel queue worker
; Save to: /etc/supervisor/conf.d/web-smk1kolaka-queue.conf

[program:web-smk1kolaka-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/web-smk1kolaka/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/web-smk1kolaka/storage/logs/queue-worker.log
stopwaitsecs=3600
environment=APP_ENV="production",QUEUE_CONNECTION="sync"
